# server/supervisor/development.ini
[program:{{ project_name }}-development]
directory = /home/thedevel/sites/{{ project_name }}-dev-env/{{ project_name }}
command = /home/thedevel/sites/{{ project_name }}-dev-env/bin/uwsgi --ini /home/thedevel/sites/{{ project_name }}-dev-env/{{ project_name }}/server/uwsgi/development.ini
autostart = true
autorestart = true
