# server/supervisor/development.ini
[program:{{ project_name }}-dev]
directory = /home/thedevel/sites/{{ project_name }}-dev-env/{{ project_name }}
command = uwsgi --ini server/uwsgi/development.ini
autostart = true
autorestart = true
