# server/supervisor/local.ini
[program:{{ project_name }}]
directory = /path/to/project
command = uwsgi --ini server/uwsgi/production.ini
autostart = true
autorestart = true
