# server/supervisor/local.ini
[program:{{ project_name }}]
directory = /home/thedevel/.virtualenvs/{{ project_name }}-env/{{ project_name }}
command = uwsgi --ini server/uwsgi/local.ini
autostart = true
autorestart = true
