# server/supervisor/local.ini
[program:{{ project_name }}-local]
directory = /home/thedevel/.virtualenvs/{{ project_name }}-env/{{ project_name }}
command = /home/thedevel/.virtualenvs/{{ project_name }}-env/bin/uwsgi --ini /home/thedevel/.virtualenvs/{{ project_name }}-env/{{ project_name }}/server/uwsgi/local.ini
autostart = true
autorestart = true
