# server/supervisor/production.ini
[program:{{ project_name }}-production]
directory = /home/thedevel/sites/{{ project_name }}-env/{{ project_name }}
command = /home/thedevel/sites/{{ project_name }}-env/bin/uwsgi --ini /home/thedevel/sites/{{ project_name }}-env/{{ project_name }}/server/uwsgi/production.ini
autostart = true
autorestart = true
