# server/supervisor/production.ini
[program:app]
command = uwsgi --ini /home/devuser/webapps/app-env/app/server/uwsgi/production.ini
autostart = true
autorestart = true
