# server/supervisor/production.ini
[program:{{ project_name }}]
directory = /home/thedevel/sites/{{ project_name }}-env/{{ project_name }}
command = uwsgi --ini server/uwsgi/production.ini
autostart = true
autorestart = true
